# templates/x.service.j2
[Unit]
Description=startx

[Service]
User=test
ExecStart=/bin/sh -c "sudo -u {{ kiosk_user }} startx /etc/X11/Xsession"

[Install]
WantedBy=multi-user.target
