# templates/x.service.j2
[Unit]
Description=startx

[Service]
User={{ kiosk_user }}
ExecStart=/bin/sh -c "sudo -S -u {{ kiosk_user }} startx /etc/X11/Xsession"

[Install]
WantedBy=multi-user.target
